# etc/systemd/system/apache2.service.d/apache2-preflight.conf
[Service]
ExecStartPre=/usr/local/scratches/bin/apache2-preflight
