# etc/cron.d/monit-checkup
# Check MONIT is running
*/2 * * * *	root	/usr/local/bin/monit-checkup
