# etc/profile.d/systemctl@user.sh
#!/bin/sh
systemctl --user status fstab-user.service
